# ogr2ogr/refresh-data.cron
0 * * * * root /bin/bash /opt/refresh-scripts/refresh-data.sh /opt/refresh-scripts/refresh-data.ini